# deploy/galle-update.service
[Unit]
Description=Pull latest Galle Font Terminal and restart kiosk if changed
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/mooniak-pi4/galle-font-terminal/deploy/galle-update.sh
